# docker/s6-overlay/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

exec php-fpm7
