# docker/s6-overlay/etc/services.d/crond/run
#!/usr/bin/with-contenv bash

exec crond -f
