# docker/s6-overlay/etc/services.d/caddy/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid caddy caddy run -config /etc/caddy/Caddyfile
